# remind.cron
29 20 9 10 4 /home/ashin/projects/Date-Remind/.venv/bin/python /home/ashin/projects/Date-Remind/send_remind.py ashin.chen@bosondata.com.cn 213 r_id46
